# server/kokoro_tts_service/tji-kokoro-tts.service
[Unit]
Description=TJI Kokoro TTS Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/tji/kokoro-tts/server
Environment=KOKORO_MODEL_DIR=/opt/tji/kokoro-tts/kokoro-multi-lang-v1_0
Environment=KOKORO_NUM_THREADS=4
Environment=KOKORO_CACHE_MAX_ITEMS=64
ExecStart=/opt/tji/kokoro-tts/venv/bin/uvicorn app:app --host 0.0.0.0 --port 8008
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
