# server/kokoro_tts_service/tji-kokoro-tts.service
[Unit]
Description=TJI Speaker Record Transfer Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/tji/kokoro-tts/server
ExecStart=/opt/tji/kokoro-tts/venv/bin/uvicorn app:app --host 0.0.0.0 --port 8008
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
